# lab-orchestrator/deploy/kr.ceml.lab-orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>kr.ceml.lab-orchestrator</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source ~/.zshrc 2>/dev/null; conda activate lab-research-agents &amp;&amp; cd ~/Dropbox/Dev/CEML_RA/lab-orchestrator &amp;&amp; python3 -m api.server</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mersoom/Dropbox/Dev/CEML_RA/lab-orchestrator/logs/orchestrator.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mersoom/Dropbox/Dev/CEML_RA/lab-orchestrator/logs/orchestrator.error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/mersoom/Dropbox/Dev/CEML_RA/lab-orchestrator</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
