# lab-orchestrator/deploy/kr.ceml.lab-ui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>kr.ceml.lab-ui</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/npm</string>
        <string>run</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mersoom/Dropbox/Dev/CEML_RA/lab-orchestrator/ui</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mersoom/Dropbox/Dev/CEML_RA/lab-orchestrator/logs/ui.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mersoom/Dropbox/Dev/CEML_RA/lab-orchestrator/logs/ui.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PORT</key>
        <string>3000</string>
    </dict>
</dict>
</plist>
