# lab-paper-scout/com.ceml.paper-scout.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Template only: copy to com.ceml.paper-scout.local.plist or ~/Library/LaunchAgents/ and fill host-local values there. -->
    <key>Label</key>
    <string>com.ceml.paper-scout</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mersoom/miniconda3/envs/lab-research-agents/bin/python</string>
        <string>/Users/mersoom/Library/CloudStorage/Dropbox/Dev/CEML_RA/lab-paper-scout/run.py</string>
        <string>daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mersoom/Library/CloudStorage/Dropbox/Dev/CEML_RA/lab-paper-scout</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>GEMINI_API_KEY</key>
        <string>REPLACE_WITH_GEMINI_API_KEY</string>
        <key>S2_API_KEY</key>
        <string>REPLACE_WITH_S2_API_KEY_OR_LEAVE_EMPTY</string>
        <key>SLACK_WEBHOOK_API</key>
        <string>REPLACE_WITH_SLACK_WEBHOOK_URL</string>
        <key>OPENAI_API_KEY</key>
        <string>REPLACE_WITH_OPENAI_API_KEY</string>
        <key>PATH</key>
        <string>/Users/mersoom/miniconda3/envs/lab-research-agents/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/mersoom/Library/CloudStorage/Dropbox/Dev/CEML_RA/lab-paper-scout/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mersoom/Library/CloudStorage/Dropbox/Dev/CEML_RA/lab-paper-scout/logs/launchd_stderr.log</string>
</dict>
</plist>
